# elasticsearch.ini
[program:elasticsearch]
command = /home/admin/elasticsearch/bin/elasticsearch 
autostart = true     ; 在 supervisord 启动的时候也自动启动
startsecs = 5        ; 启动 5 秒后没有异常退出，就当作已经正常启动了
autorestart = true   ; 程序异常退出后自动重启
user = admin 
redirect_stderr = true  
stdout_logfile_maxbytes = 20MB  
stdout_logfile_backups = 20    
stdout_logfile = /data/elasticsearch/logs/elasticsearch.log
