# elasticsearch.ini
[unix_http_server]
file=/home/admin/elasticsearch/run/supervisor_es.sock

[supervisord]
directory=/home/admin/elasticsearch
logfile=./logs/elasticsearch.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=warn
pidfile=./run/supervisor_elasticsearch.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/admin/elasticsearch/run/supervisor_es.sock

[program:elasticsearch]
stopsignal=INT
directory=/home/admin/elasticsearch
command=/home/admin/elasticsearch/bin/elasticsearch
user=admin
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/data/elasticsearch/logs/elasticsearch.log
stderr_logfile=/data/elasticsearch/logs/elasticsearch.log

