# kafka.ini
[unix_http_server]
file=/home/admin/kafka/run/supervisor_kafka.sock

[supervisord]
directory=/home/admin/kafka
logfile=./logs/kafka.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=debug
pidfile=./run/supervisord_kafka.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/admin/kafka/run/supervisor_kafka.sock

[program:kafka]
stopsignal=INT
directory=/home/admin/kafka
command=/home/admin/kafka/bin/kafka-server-start.sh /home/admin/kafka/config/server.properties
user=admin
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/home/admin/kafka/logs/kafka.log
stderr_logfile=/home/admin/kafka/logs/kafka.log
