# kafka.ini
[program:kafka]
command=/home/admin/kafka/bin/kafka-server-start.sh /home/admin/kafka/config/server.properties
autostart=true
autorestart=true
startsecs=3
user=admin
redirect_stderr=true
stdout_logfile=/home/admin/kafka/logs/kafka.log
