# redis.ini
[program:redis]
command=/usr/local/bin/redis-server /home/admin/redis/conf/redis.conf
stdout_logfile=/home/admin/redis6379/logs/redis.log
stdout_logfile_maxbytes=9MB
stdout_logfile_backups=10
user=admin
autorestart=true ;自动重启
