# redis.ini
[unix_http_server]
file=/home/admin/redis/run/supervisor_redis.sock

[supervisord]
directory=/home/admin/redis
logfile=./logs/redis_supervisor.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=warn
pidfile=./run/supervisor_redis.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:////home/admin/redis/run/supervisor_redis.sock

[program:redis]
stopsignal=INT
directory=/home/admin/redis
command=./sbin/redis-server ./conf/redis.conf
user=admin
autostart=true
autorestart=true
stdout_logfile=/home/admin/redis/logs/redis.log
stderr_logfile=/home/admin/redis/logs/redis.log
