# zookeeper.ini
[unix_http_server]
file=/home/admin/zookeeper/run/supervisor_zookeeper.sock

[supervisord]
directory=/home/admin/zookeeper
logfile=./logs/zookeeper_supervisor.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=warn
pidfile=./run/supervisor_zookeeper.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:////home/admin/zookeeper/run/supervisor_zookeeper.sock

[program:zookeeper]
stopsignal=INT
directory=/home/admin/zookeeper
command=./bin/zkServer.sh start-foreground
user=admin
autostart=true
autorestart=true
stdout_logfile=/home/admin/zookeeper/logs/zookeeper.log
stderr_logfile=/home/admin/zookeeper/logs/zookeeper.log
