# zookeeper.ini
[program:zookeeper]
command=/home/admin/zookeeper/bin/zkServer.sh start-foreground
stdout_logfile=/home/admin/zookeeper/logs/zookeeper.out
stdout_logfile_maxbytes=9MB
stdout_logfile_backups=10
user=admin
autorestart=true ;自动重启
